# deploy/dc-frontier-events-audit.timer
[Unit]
Description=Run the DC Frontier Events live audit weekly

[Timer]
OnCalendar=Mon 06:00
Persistent=true

[Install]
WantedBy=timers.target
